# alpine/tt-system-tools.pre-deinstall
#!/bin/sh

rc-update delete tenstorrent-hugepages 
rc-update delete tenstorrent-mount-hugepages 